# agents/oleg-mac/com.daily.inject.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.daily.inject</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>/Users/user/daily-inject.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>600</integer>

    <key>StandardOutPath</key>
    <string>/tmp/daily-inject.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/daily-inject.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>GIT_SSH_COMMAND</key>
        <string>ssh -i /Users/user/.ssh/github_ed25519 -o StrictHostKeyChecking=no</string>
    </dict>
</dict>
</plist>
